# debian/postinst
#!/bin/sh
set -e

# Create the .local/share directory if it doesn't exist
if [ ! -d ~/.local/share ]; then
    mkdir -p ~/.local/share
fi

# Enable systemd user service for all users
if [ -f /usr/lib/systemd/user/reminder.service ]; then
    echo "Enabling systemd user service for all users"
    # This will only work for currently logged-in users
    # New users will need to enable it manually or through user session startup
    for uid in $(loginctl list-users | awk 'NR>1 {print $1}'); do
        username=$(loginctl show-user $uid -p Name | cut -d= -f2)
        if [ -n "$username" ] && [ "$username" != "root" ]; then
            homedir=$(getent passwd $username | cut -d: -f6)
            if [ -d "$homedir" ]; then
                echo "Enabling service for user $username"
                # Use su to run the command as the user
                su - $username -c "systemctl --user daemon-reload && systemctl --user enable reminder.service"
            fi
        fi
    done
fi

#DEBHELPER#

exit 0
